# templates/st2api.socket.j2
[Unit]
Description=StackStorm st2api Socket
PartOf=st2api.service

[Socket]
ListenStream={{ ansible_default_ipv4.address }}:{{ st2_port.api }}

[Install]
WantedBy=sockets.target
