# templates/st2stream.socket.j2
[Unit]
Description=StackStorm st2stream Socket
PartOf=st2stream.service

[Socket]
ListenStream=127.0.0.1:{{ st2_port.stream }}

[Install]
WantedBy=sockets.target
